# files/nfqws.service
[Unit]
Description=nfqwscfg
After=network.target

[Service]
ExecStart=/usr/bin/zapret.sh

[Install]
WantedBy=multi-user.target
